# backend/smartvenue-network-scanner.service
[Unit]
Description=SmartVenue Network Scanner Service
After=network.target

[Service]
Type=simple
User=coastal
WorkingDirectory=/home/coastal/smartvenue/backend
ExecStart=/home/coastal/smartvenue/backend/network_scanner_scheduled.sh --interval 10 --subnet 192.168.101
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
